# deploy/flowlens-agent.service
[Unit]
Description=FlowLens host traffic agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=flowlens-agent
Group=flowlens-agent
EnvironmentFile=/etc/flowlens/agent.env
ExecStart=/usr/local/bin/flowlens-agent
Restart=on-failure
RestartSec=3s
TimeoutStopSec=10s
MemoryMax=96M
TasksMax=128
UMask=0027
NoNewPrivileges=true
AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW
CapabilityBoundingSet=CAP_BPF CAP_PERFMON CAP_SYS_RESOURCE CAP_NET_RAW
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_PACKET
ReadWritePaths=/var/lib/flowlens-agent

[Install]
WantedBy=multi-user.target
