# deploy/flowlens-server.service
[Unit]
Description=FlowLens server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=flowlens
Group=flowlens
EnvironmentFile=/etc/flowlens/server.env
ExecStart=/usr/local/bin/flowlens-server
Restart=on-failure
RestartSec=3s
TimeoutStopSec=15s
MemoryMax=384M
TasksMax=128
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/flowlens

[Install]
WantedBy=multi-user.target
